# deploy/warmdesk.service
[Unit]
Description=WarmDesk Kanban
Documentation=https://github.com/tonk/warmdesk
After=network.target

[Service]
Type=simple
User=warmdesk
Group=warmdesk
WorkingDirectory=/opt/warmdesk

ExecStart=/opt/warmdesk/warmdesk
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Environment — adjust to your setup
Environment=PORT=8080
Environment=DB_DRIVER=sqlite
Environment=DB_DSN=/opt/warmdesk/data/warmdesk.db
Environment=JWT_SECRET=change-me-to-a-random-secret
Environment=ALLOWED_ORIGINS=https://yourdomain.com
Environment=WEB_DIR=/opt/warmdesk/web

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/warmdesk/data

[Install]
WantedBy=multi-user.target
